# venv/routerconfig.ini
[program:yxtconfig]
heartbeat_time = 600
order_time = 300
at_time = 120
server_url = http://wifi.jzlkbj.com:8181/api/
ver_code = 18
on_time = 20180809110123
retry_time = 3
